# tor/hidden_service/honeypot/hidden-service/contrib/supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
command=/docker-entrypoint.sh tor

[program:tap]
command=/usr/local/bin/tap.sh
