# tor/hidden_service/honeypot/hidden-service/contrib/supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
command=/docker-entrypoint.sh tor

[program:tap]
command=tcpdump -i eth0 -s 0 -w /dev/stdout | nc bro 1969

