# tor/hidden_service/hidden-service/assets/supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
command=/docker-entrypoint.sh tor
